# backend/systemd/ollama-user.service
# Ollama user service (no sudo required).
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp ollama-user.service ~/.config/systemd/user/ollama.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now ollama
#
# Optional — start at boot without logging in (needs sudo once):
#   sudo loginctl enable-linger $USER
#
# For system-wide install (starts before Docker), use ollama.service instead:
#   sudo cp ollama.service /etc/systemd/system/
#   sudo systemctl enable --now ollama

[Unit]
Description=Ollama LLM server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ollama serve
Restart=on-failure
RestartSec=5

Environment=OLLAMA_HOST=0.0.0.0:11434
Environment=OLLAMA_KEEP_ALIVE=15m
Environment=OLLAMA_CONTEXT_LENGTH=2048
Environment=OLLAMA_DEBUG=INFO
Environment=HOME=%h

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama

[Install]
WantedBy=default.target
